# local/com.buttonup.bumbot-review-text.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the local iMessage client-review reminder. Reinstall after edits:
    cp local/com.buttonup.bumbot-review-text.plist ~/Library/LaunchAgents/
    launchctl unload ~/Library/LaunchAgents/com.buttonup.bumbot-review-text.plist 2>/dev/null
    launchctl load -w ~/Library/LaunchAgents/com.buttonup.bumbot-review-text.plist

  CATCH-UP: it runs at 11:00 (StartCalendarInterval), on login/boot (RunAtLoad),
  on wake, and hourly (StartInterval). The SCRIPT decides whether to actually send
  (weekday, ≥11 AM ET, not already sent today), so a reminder missed because the
  Mac was asleep or powered off fires the next time the Mac is awake — once.

  First send may prompt to allow osascript to control Messages
  (System Settings → Privacy & Security → Automation). Logs: /tmp/bumbot-review-text.*
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.buttonup.bumbot-review-text</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/christianbusiness/BUMBOT Master Sheet/local/reviewTextAgent.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/christianbusiness/BUMBOT Master Sheet</string>
  <!-- Run on login/boot and immediately on (re)load — catches a missed 11am after a shutdown. -->
  <key>RunAtLoad</key>
  <true/>
  <!-- Hourly safety net; also fires on wake if the interval elapsed while asleep. -->
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/bumbot-review-text.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bumbot-review-text.err</string>
</dict>
</plist>
