# local/com.buttonup.bumbot-review-text.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: runs the local iMessage client-review reminder weekdays at 11:00
  (Mac LOCAL time). Install:
    cp local/com.buttonup.bumbot-review-text.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.buttonup.bumbot-review-text.plist
  (or `launchctl load -w ~/Library/LaunchAgents/com.buttonup.bumbot-review-text.plist` on older macOS).
  It only fires when the Mac is awake at 11am. First run may prompt to allow
  Terminal/osascript to control Messages (System Settings → Privacy → Automation).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.buttonup.bumbot-review-text</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/christianbusiness/BUMBOT Master Sheet/local/reviewTextAgent.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/christianbusiness/BUMBOT Master Sheet</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/bumbot-review-text.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bumbot-review-text.err</string>
</dict>
</plist>
